# etc/crontabs/root
# m h dom mon dow command
# Run homepass every 2 minutes
*/2 *  *   *   *  /usr/bin/homepass.sh 2>&1 >> /tmp/homepass.log
# Clean log file every 8 hours
 0 */8 *   *   *  echo "$(date): Log truncated." > /tmp/homepass.log
